# launchd/com.worldfacts.commandbot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.worldfacts.commandbot</string>

    <key>ProgramArguments</key>
    <array>
      <string>/Users/antonzaritskyi/Desktop/VS studio Code/project2/.venv/bin/python</string>
      <string>/Users/antonzaritskyi/Desktop/VS studio Code/project2/telegram_command_bot.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/antonzaritskyi/Desktop/VS studio Code/project2</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/antonzaritskyi/Desktop/VS studio Code/project2/logs/commandbot.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/antonzaritskyi/Desktop/VS studio Code/project2/logs/commandbot.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
      <key>PYTHONUNBUFFERED</key>
      <string>1</string>
    </dict>
  </dict>
</plist>
